# fire-department-bp-standalone/bp-management.service
[Unit]
Description=Fire Department BP Management System
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/bp-management
ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
